# scripts/online-judge.service
[Unit]
Description=Online Judge Next.js App
After=network.target
OnFailure=notify-failure@%n.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/online-judge
Environment=NODE_ENV=production
EnvironmentFile=/home/ubuntu/online-judge/.env
ExecStart=/usr/bin/npm start -- -H 127.0.0.1 -p 3000
Restart=always
RestartSec=5
# Security hardening
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
NoNewPrivileges=true
ReadWritePaths=/home/ubuntu/online-judge/data
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
MemoryMax=2G
CPUQuota=200%

[Install]
WantedBy=multi-user.target
