# scripts/online-judge-worker.service
[Unit]
Description=Online Judge Judge Worker
After=network.target docker.service online-judge.service
Requires=online-judge.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/online-judge
EnvironmentFile=/home/ubuntu/online-judge/.env
Environment=NODE_ENV=production
ExecStart=/usr/bin/npx tsx judge-worker/index.ts
Restart=always
RestartSec=5
# Security hardening
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
NoNewPrivileges=true
ReadWritePaths=/home/ubuntu/online-judge/data
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
CapabilityBoundingSet=
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
